# debian/postinst
#!/bin/bash
# Utappia repository installer
# Copyright (c) 2017 Salih Emin. All rights reserved.
# Use of this source code is governed by GPLv3 license that can be
# found in the copyright file.
set -e

function add_debian_repository() {
  echo "** Adding Utappia public key"
  apt-key adv --keyserver keyserver.ubuntu.com --recv EBC8F9FD
  echo "** Adding Debian repository"
  cat <<EOT >/etc/apt/sources.list.d/utappia.list
# Debian compatible Utappia Repository
deb http://ppa.launchpad.net/utappia/stable/ubuntu xenial main
deb-src http://ppa.launchpad.net/utappia/stable/ubuntu xenial main

EOT
}

echo "** Checking if Utappia Repository exists in your system"
if ls /etc/apt/sources.list.d/utappia* 1>/dev/null 2>&1; then
  echo "Utappia Repository exists... skiping"
else
  if ! [ -x "$(command -v add-apt-repository)" ]; then
    echo "** Adding Utappia repository"
    add-apt-repository -y ppa:utappia/stable
  else
    add_debian_repository
  fi
fi

#if [[ -x "$(command -v  systemctl )" ]]; then
#  echo "** Enable systemd services"
#  systemctl enable ucaresystem-automation-cleanup.service || true
#fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
